# custom_metrics.target
[Unit]
Description=Custom Metrics Service Target
StopWhenUnneeded=no