# custom_metrics.timer
[Unit]
Description=Minute Timer

[Timer]
OnBootSec=1min
OnCalendar=*:0/1
Unit=custom_metrics.target

[Install]
WantedBy=custom_metrics.target
